# configs/systemd/benchmark.slice
[Unit]
Description=Benchmark workload slice
Before=slices.target

[Slice]
